# runner/launchd/ai.agents.notion.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.agents.notion.scheduler</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/matheuszimmermannai/agents/runner/notion_scheduler.py</string>
    </array>

    <!-- roda a cada 60 minutos -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>SSL_CERT_FILE</key>
      <string>/usr/local/lib/python3.14/site-packages/certifi/cacert.pem</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/matheuszimmermannai/agents/runner/logs/notion.scheduler.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/matheuszimmermannai/agents/runner/logs/notion.scheduler.err</string>
  </dict>
</plist>
