# runner/launchd/ai.agents.posts.secureapix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.agents.runner</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/matheuszimmermannai/Documents/agents/runner/run_jobs.py</string>
    </array>

    <!-- roda todo dia às 09:00 -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>

    <!-- roda também quando você carregar / quando o Mac ligar (catch-up) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/matheuszimmermannai/Documents/agents/runner/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/matheuszimmermannai/Documents/agents/runner/logs/launchd.err</string>
  </dict>
</plist>
